# kual-extension/kindle-series-manager/upstart/gr-sync.conf
description "Goodreads progress sync daemon"

start on framework_ready and started lab126_gui
stop on stopping system

export LANG=en_US.UTF-8

script
    LOG="/mnt/us/gr-upstart.log"
    EXT_DIR="/mnt/base-us/extensions/kindle-series-manager"
    FLAG_FILE="/mnt/us/ENABLE_GR_SYNC"
    SYNC_SCRIPT="$EXT_DIR/goodreads/gr_sync.sh"

    echo "$(date) gr-sync starting" >> "$LOG"

    if [ ! -f "$FLAG_FILE" ]; then
        echo "$(date) FLAG_FILE missing, exiting" >> "$LOG"
        exit 0
    fi

    if [ ! -f "$SYNC_SCRIPT" ]; then
        echo "$(date) SYNC_SCRIPT missing, exiting" >> "$LOG"
        exit 0
    fi

    echo "$(date) launching gr_sync.sh" >> "$LOG"
    exec bash "$SYNC_SCRIPT" >> "$LOG" 2>&1
end script
