# kual-extension/kindle-series-manager/upstart/ksm-stats.conf
description "KSM reading statistics daemon"

start on framework_ready and started lab126_gui
stop on stopping system

export LANG=en_US.UTF-8

script
    LOG="/mnt/us/ksm-stats-upstart.log"
    EXT_DIR="/mnt/us/extensions/kindle-series-manager"
    FLAG_FILE="/mnt/us/ENABLE_KSM_STATS"
    DAEMON_SCRIPT="$EXT_DIR/bin/ksm_stats_daemon.sh"

    echo "$(date) ksm-stats starting" >> "$LOG"

    if [ ! -f "$FLAG_FILE" ]; then
        echo "$(date) FLAG_FILE missing, exiting" >> "$LOG"
        exit 0
    fi

    if [ ! -f "$DAEMON_SCRIPT" ]; then
        echo "$(date) DAEMON_SCRIPT missing, exiting" >> "$LOG"
        exit 0
    fi

    echo "$(date) launching ksm_stats_daemon.sh" >> "$LOG"
    exec sh "$DAEMON_SCRIPT" >> "$LOG" 2>&1
end script
